# config/launchagents/com.urfael.dashboard.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- Urfael dashboard — token-gated LOCALHOST web console (binds 127.0.0.1 ONLY, never the LAN). install.sh fills
     {{HOME}}/{{NODE}}/{{REPO}}. The page token is generated on first start and written 0600 to
     ~/.claude/urfael/dashboard.token; run `urfael dashboard` to print the tokened URL. Read SECURITY.md first. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.urfael.dashboard</string>
  <key>ProgramArguments</key>
  <array><string>{{NODE}}</string><string>{{REPO}}/app/dashboard.js</string></array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><dict><key>SuccessfulExit</key><false/></dict>
  <key>Umask</key><string>0077</string><!-- StandardOut/Err logs created 0600, never world-readable -->

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- Loopback port for the dashboard (default 7717): -->
    <!-- <key>URFAEL_DASHBOARD_PORT</key><string>7717</string> -->
  </dict>
  <key>StandardOutPath</key><string>{{HOME}}/.claude/urfael/dashboard.out.log</string>
  <key>StandardErrorPath</key><string>{{HOME}}/.claude/urfael/dashboard.err.log</string>
</dict></plist>
